# cassandra/roles/java/templates/etc/profile.d/java.sh.j2
export JAVA_HOME={{java_home}}
{% if java_jdk_or_jre == "jre" %}
export JRE_HOME={{java_home}}
{% endif %}
